# units/emacs.service
[Unit]
Description=Emacs

[Service]
EnvironmentFile=%h/.gpg-agent-info
ExecStart=/usr/bin/emacs

[Install]
WantedBy=default.target
