# units/firefox.service
[Unit]
Description=Firefox
After=xorg.target

[Service]
ExecStart=/usr/bin/firefox

[Install]
WantedBy=default.target
